# provisioning/templates/etc/init/kifi-node-setup.conf.j2
# {{ ansible_managed }}
start on started networking
task
script
  exec /opt/kifi-node-setup
end script
