# provisioning/templates/etc/init/kifi-register-host.conf.j2
# {{ ansible_managed }}
start on started networking
task
script
  exec /opt/kifi-register-host
end script
